# voice_mode/templates/launchd/com.voicemode.livekit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- com.voicemode.livekit.plist v1.0.0 -->
<!-- Last updated: 2025-08-12 -->
<!-- Compatible with: LiveKit v1.7.0+ -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicemode.livekit</string>
    <key>ProgramArguments</key>
    <array>
        <string>{LIVEKIT_BIN}</string>
        <string>--config</string>
        <string>{CONFIG_FILE}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{LOG_DIR}/livekit.out.log</string>
    <key>StandardErrorPath</key>
    <string>{LOG_DIR}/livekit.err.log</string>
    <key>WorkingDirectory</key>
    <string>{WORKING_DIR}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>LIVEKIT_URL</key>
        <string>ws://localhost:{LIVEKIT_PORT}</string>
        <key>LIVEKIT_API_KEY</key>
        <string>devkey</string>
        <key>LIVEKIT_API_SECRET</key>
        <string>secret</string>
    </dict>
</dict>
</plist>