# voice_mode/templates/launchd/com.voicemode.whisper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- com.voicemode.whisper.plist v1.1.0 -->
<!-- Last updated: 2025-08-25 -->
<!-- Uses unified startup script for dynamic model selection -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicemode.whisper</string>
    <key>ProgramArguments</key>
    <array>
        <string>~/.voicemode/services/whisper/bin/start-whisper-server.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/.voicemode/logs/whisper/whisper.out.log</string>
    <key>StandardErrorPath</key>
    <string>~/.voicemode/logs/whisper/whisper.err.log</string>
    <key>WorkingDirectory</key>
    <string>~/.voicemode/services/whisper</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>